# chaos_calmer/target/linux/mvebu/base-files/etc/hotplug.d/usb/10-usb
#!/bin/sh

# Copyright (C) 2009 OpenWrt.org

case "$PRODUCT" in
        '1d6b'*)
		# Ignore base USB hardware
        ;;
        *)
		case "$ACTION" in
			add)
				# update LEDs
				case "$DEVPATH" in
					*'usb1'*)
						# USB 1.1 or 2.0 device added to USB 1
						/usr/sbin/ledctrl 5 on
					;;
					*'usb3'*)
						# USB 1.1 or 2.0 device added to USB 2
						/usr/sbin/ledctrl 6 on
					;;
					*'usb4'*)
						# USB 3.0 device added to USB 2
						/usr/sbin/ledctrl 6 on
						/usr/sbin/ledctrl 7 on
					;;
				esac
			;;
			remove)
				# update LEDs
				case "$DEVPATH" in
					*'usb1'*)
						# USB 1.1 or 2.0 device removed from USB 1
						/usr/sbin/ledctrl 5 off
					;;
					*'usb3'*)
						# USB 1.1 or 2.0 device removed from USB 2
						/usr/sbin/ledctrl 6 off
					;;
					*'usb4'*)
						# USB 3.0 device removed from USB 2
						/usr/sbin/ledctrl 6 off
						/usr/sbin/ledctrl 7 off
					;;
				esac
			;;
		esac
	;;
esac
